# server/supervisord.conf
[program:node]
command = node server.js 0.0.0.0 8080
directory = /home/dotcloud/current
