# server/supervisord.conf
[program:node]
command = node server.js
directory = /home/dotcloud/current
